# ros2_ws/src/XRF2/XenoRosFramework/Common/xbuf_dev_rule/99-xbuf-permissions.rules
# Set permissions for files in /dev/evl/xbuf to 666
ACTION=="add", SUBSYSTEM=="xbuf", DEVPATH=="/devices/virtual/xbuf/*", RUN+="/bin/chmod 666 /dev/evl/xbuf/%k"
